# distro/config/includes.chroot/etc/systemd/system/cyjan-tap-disk-watch.service
[Unit]
Description=Auto-Prune Kafka topics on Tap when disk > 85%
Documentation=https://github.com/JxxKal/ids
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
# Check disk usage and trigger cyjan-tap prune if /-Mount > 85%.
# Idempotent: bei <85% no-op.
ExecStart=/bin/bash -c '\
  USAGE=$(df / | awk "NR==2 {gsub(\"%\",\"\"); print $5}"); \
  if [ "$USAGE" -ge 85 ]; then \
    echo "[disk-watch] / at ${USAGE}% — Triggering cyjan-tap prune"; \
    /usr/local/bin/cyjan-tap prune; \
  else \
    echo "[disk-watch] / at ${USAGE}% — OK"; \
  fi'
