# distro/config/includes.chroot/etc/systemd/system/cyjan-tap-disk-watch.service
[Unit]
Description=Auto-Prune Kafka topics on Tap when disk > 85%
Documentation=https://github.com/JxxKal/ids
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
# Logik liegt im Skript /usr/local/bin/cyjan-tap-disk-watch (Tap-only,
# self-guarding über /etc/cyjan/mode). Standalone statt inline-ExecStart, weil
# die awk-in-bash-in-systemd-Quoting-Verschachtelung `$5` verschluckte und
# `awk: syntax error at or near %` warf. Idempotent: <85% = no-op.
ExecStart=/usr/local/bin/cyjan-tap-disk-watch
