# distro/config/includes.chroot/etc/systemd/system/cyjan-tap-update.service
[Unit]
Description=Cyjan Tap — Update vom Master ziehen (host-side)
Documentation=https://github.com/JxxKal/ids
After=docker.service
Requires=docker.service
ConditionPathExists=/usr/local/bin/cyjan-tap

[Service]
Type=oneshot
# --auto: nicht-interaktiv (der Master-Trigger gilt als Confirmation), aber
# OHNE Force — ist die angebotene Version identisch zur installierten, beendet
# sich cyjan-tap sauber statt einen erneuten Stack-Recreate auszulösen. Ein
# per PathChanged doppelt gefeuerter Trigger kann so keine Re-Update-Schleife
# starten. (Früher -y → bedingungsloses Re-Apply auch bei gleicher Version.)
ExecStart=/usr/local/bin/cyjan-tap update --from-master --auto
# Nach erfolgreichem Update das Trigger-File wegräumen — sonst feuert
# der path-Watcher beim nächsten ähnlichen Schreib-Event nochmal. Das
# ExecStartPost läuft NUR bei rc=0.
ExecStartPost=/bin/rm -f /run/cyjan-update/trigger
# Failsoft: wenn der Update fehlschlägt, nicht den path-Watcher selbst
# in Endlosschleife schicken. Trigger-File bleibt erhalten, beim
# nächsten manuellen Touch oder Reboot wird es nochmal versucht.
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600
