# distro/config/includes.chroot/etc/systemd/system/cyjan-stack-alert.service
[Unit]
Description=Cyjan IDS Stack-Alarm (Boot-Health-Check fehlgeschlagen)
Documentation=https://github.com/JxxKal/ids

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cyjan-stack-alert
StandardOutput=journal
StandardError=journal
